# pyuftpd/build-tools/src/main/package/distributions/Debian/debian/unicore-uftpd.service
[Unit]
Description=UNICORE UFTPD

[Service]
Type=forking
ExecStart=/usr/share/unicore/uftpd/bin/unicore-uftpd-start.sh
ExecStop=/usr/share/unicore/uftpd/bin/unicore-uftpd-stop.sh
PIDFile=/var/run/unicore-uftpd.pid

[Install]
#WantedBy=multi-user.target
